# conf/bitfinex-init.conf
start on startup
chdir /root/bitcoinelasticsearch/python-scripts
exec python3 bitfinex_websockets_connector.py
respawn