# conf/init/connector-init.conf
start on runlevel [23456]
chdir /home/ubuntu/bitcoinelasticsearch/python-scripts
exec python3 __PYTHONFILE__ --host=__ESHOST__
respawn
