# conf/init/okcoin-connector-init.conf
start on runlevel [23456]
chdir /home/ubuntu/bitcoinelasticsearch/src/python/okcoin
exec python3 okcoin-wrapper.py
respawn
