# conf/init/okcoin-init.conf
start on startup
chdir /root/bitcoinelasticsearch/python-scripts
exec python3 okcoin_websockets_connector.py
respawn