# conf/init/okcoin-init.conf
start on runlevel [23456]
chdir /root/bitcoinelasticsearch/python-scripts
exec python3 okcoin_websockets_connector.py
respawn